# benchmarking_sbi/ansible/templates/rq_workers_supervisord.conf.j2
[supervisord]
; http://supervisord.org/configuration.html#supervisord-section-values


[supervisorctl]


[inet_http_server]
port = 0.0.0.0:9001


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:rq_worker]
environment=PATH='{{ rq_workers_bin }}:%(ENV_PATH)s',MKL_NUM_THREADS=1,NUMEXPR_NUM_THREADS=1,OMP_NUM_THREADS=1,OPENBLAS_NUM_THREADS=1
command={{ rq_workers_bin }}/rqworker {{ rq_workers_queue }} --url {{ rq_workers_redis_url }}
directory=%(here)s
numprocs=%(ENV_NUM_WORKERS)s
process_name=%(program_name)s-%(process_num)s
stopsignal=TERM
autorestart=true
autostart=true
